# bin/s6-rc.d/hadoop/run
#!/bin/sh
/usr/sbin/sshd
HOSTNAME=$(hostname)
source /etc/profile
if [[ "$HOSTNAME" == *"slave1"* ]]; then
    echo "1" > /usr/local/software/zookeeper-3.9.0/datadir/myid
fi
if [[ "$HOSTNAME" == *"slave2"* ]]; then
    echo "2" > /usr/local/software/zookeeper-3.9.0/datadir/myid
fi
if [[ "$HOSTNAME" == *"master"* ]]; then
    echo "0" > /usr/local/software/zookeeper-3.9.0/datadir/myid
    /usr/local/software/hadoop-3.3.6/sbin/start-all.sh
    /usr/local/software/hbase-2.5.5/bin/start-hbase.sh
    /usr/local/software/hadoop-3.3.6/bin/hadoop fs -mkdir -p /spark/spark-log
    /usr/local/software/spark-3.4.1/sbin/start-all.sh
    /usr/local/software/spark-3.4.1/sbin/start-history-server.sh
    /root/bin/zk start
    echo "Y" | yarn rmadmin -transitionToStandby --forcemanual rm2
    echo "Y" | yarn rmadmin -transitionToActive --forcemanual rm1
    MYSQL_USER="root"
    MYSQL_PASSWORD="root"
    MYSQL_HOST="mysql"
    MYSQL_DATABASE="hive"
    # 验证数据库是否存在
    if mysql -u $MYSQL_USER -p$MYSQL_PASSWORD -h $MYSQL_HOST -e "use $MYSQL_DATABASE"; then
    echo "数据库已存在，跳过初始化"
    else
    /usr/local/software/hive-3.1.3/bin/schematool -initSchema -dbType mysql
    echo "数据库初始化完成"
    fi
    nohup /usr/local/software/hive-3.1.3/bin/hive --service metastore 2>&1 &
    echo "启动hive元数据服务"
    echo "等待启动HiveServer2服务"
    while true; do
        # 检查端口是否开放
        if netstat -an | grep 9083 > /dev/null; then
            # 启动HiveServer2服务
            /usr/local/software/hadoop-3.3.6/bin/hdfs dfsadmin -safemode leave
            nohup /usr/local/software/hive-3.1.3/bin/hive --service hiveserver2 2>&1 &
            echo "启动HiveServer2服务"
            break
        else
            echo "hive元数据服务未就绪，等待中..."
            sleep 1
        fi
    done
fi
echo "$HOSTNAME"启动完成
