# bin/s6-rc.d/hadoop/run
#!/bin/sh
/usr/sbin/sshd
HOSTNAME=$(hostname)
source /etc/profile
if [[ "$HOSTNAME" == *"slave1"* ]]; then
    echo "1" > /usr/local/software/zookeeper-3.8.3/datadir/myid
fi
if [[ "$HOSTNAME" == *"slave2"* ]]; then
    echo "2" > /usr/local/software/zookeeper-3.8.3/datadir/myid
fi
if [[ "$HOSTNAME" == *"master"* ]]; then
    echo "0" > /usr/local/software/zookeeper-3.8.3/datadir/myid
    /usr/local/software/hadoop-3.3.6/sbin/start-all.sh
    /usr/local/software/hbase-2.5.5/bin/start-hbase.sh
    /usr/local/software/hadoop-3.3.6/bin/hadoop fs -mkdir -p /spark/spark-log
    /usr/local/software/spark-3.5.0/sbin/start-all.sh
    /usr/local/software/spark-3.5.0/sbin/start-history-server.sh
    /root/bin/zk start
    echo "Y" | yarn rmadmin -transitionToStandby --forcemanual rm2 > /dev/null 2>&1
    echo "Y" | yarn rmadmin -transitionToActive --forcemanual rm1 > /dev/null 2>&1
    # 验证数据库是否存在
    if mysql -u root -proot -h mysql -e "use hive" > /dev/null 2>&1; then
        echo "数据库已存在，跳过初始化"
    else
        /usr/local/software/hive-3.1.3/bin/schematool -initSchema -dbType mysql > /dev/null 2>&1
        echo "数据库初始化完成"
    fi
    nohup /usr/local/software/hive-3.1.3/bin/hive --service metastore 2>&1 &
    echo "启动hive元数据服务"
    echo "等待启动HiveServer2服务"
    echo "hive元数据服务未就绪，等待中..."
    while true; do
        # 检查端口是否开放
        if netstat -an | grep 9083 > /dev/null; then
            # 启动HiveServer2服务
            /usr/local/software/hadoop-3.3.6/bin/hdfs dfsadmin -safemode leave
            nohup /usr/local/software/hive-3.1.3/bin/hive --service hiveserver2 2>&1 &
            echo "启动HiveServer2服务"
            break
        else
            sleep 1
        fi
    done
    while true; do
        # 尝试连接到Hive Server 2并执行查询
        beeline -u "jdbc:hive2://localhost:10000" -e "SELECT 1;" > /dev/null 2>&1
        if [ $? -eq 0 ]; then
            echo "HiveServer2服务启动完成"
            break
        else
            sleep 1
        fi
    done
fi
echo "$HOSTNAME"启动完成
